# gpu-cooling.service
[Unit]
Description=GPU Temperature Based Cooling Controller
After=multi-user.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/bin/bash /usr/local/bin/gpu-cooling-wrapper.sh

[Install]
WantedBy=multi-user.target
